# files/system/etc/systemd/system/mnt-nas-media.mount
[Unit]
Description=NFS mount for media storage (Plex/Jellyfin)
After=network-online.target
Wants=network-online.target
Before=docker.service

[Mount]
What=192.168.7.10:/mnt/storage/Media
Where=/mnt/nas-media
Type=nfs
Options=ro,hard,intr,rsize=131072,wsize=131072,tcp,timeo=600,retrans=2,_netdev

# Retry on failure
TimeoutSec=60

[Install]
WantedBy=multi-user.target
WantedBy=remote-fs.target
