# files/system/etc/profile.d/intel-vaapi.sh
#!/usr/bin/env sh
# Intel VAAPI configuration for hardware acceleration
# Sets the default VAAPI driver to iHD (intel-media-driver)
export LIBVA_DRIVER_NAME=iHD
